# oj_spider.ini
[program: oj_spider-app]
command = python3 app.py
directory = /home/oj_spider
autostart = true

[program: oj_spider-celery]
command = celery -A tasks worker -l info -c 8 --pool=eventlet
directory = /home/oj_spider
autostart = true

[program: oj_spider-beat]
command = celery -A tasks beat -l info
directory = /home/oj_spider
autostart = true

