# conf/beat.ini
[program: beat]
command = celery -A tasks beat -l info
directory = /home/oj_spider
autostart = true