# conf/celery.ini
[program: celery]
command = celery -A tasks worker -l info -c 8 --pool=eventlet
directory = /home/oj_spider
autostart = true